# klipper/mcu.ini
[program:mcu]
user=root
command=/home/klippy/klipper/out_linux/klipper.elf -w
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes = 0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes = 0
autorestart=true
autostart=false
startsecs=0
dependent_startup=true
